# v2ray/supervisor.d/v2ray.ini
[program:v2ray-server]
command=/bin/sh -c '
    DOMAIN=$(ls /etc/nginx/certs |grep 0x256 | sed -n 1p)
    if [ "${DOMAIN}" != "" ]; then
        sed "s/DOMAIN/${DOMAIN}/g" /etc/v2ray/server.tmp > /etc/v2ray/server.json
        sed "s/DOMAIN/${DOMAIN}/g" /etc/nginx/http.d/default.tmp > /etc/nginx/http.d/default.conf
        /usr/sbin/nginx -c /etc/nginx/nginx.conf
    else
        sed "6,35d;129,153d" /etc/v2ray/server.tmp > /etc/v2ray/server.json
    fi
    /usr/bin/v2ray -config /etc/v2ray/server.json'
process_name=%(program_name)s
numprocs=1
startsecs=3
autostart=true
autorestart=true
priority=1
user=root
umask=0000
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
