# v2ray/supervisor.d/v2ray.ini
[program:v2ray-server]
command=/bin/sh -c '
    sed "/DOMAIN/%(ENV_DOMAIN)s" /etc/v2ray/server.sample.json > /etc/v2ray/server.json
    /usr/bin/v2ray run -c /etc/v2ray/server.json'
process_name=%(program_name)s
numprocs=1
startsecs=3
autostart=true
autorestart=true
priority=1
user=root
umask=0000
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
