# syncthing/supervisor.d/default.ini
[supervisord]
pidfile=/var/run/supervisord.pid
nodaemon=true

[program:syncthing]
directory=%(ENV_SYNCTHING_PATH)s
command=/bin/entrypoint.sh /bin/syncthing --config=%(ENV_SYNCTHING_PATH)s --data=%(ENV_DATA_PATH)s
process_name=%(program_name)s
numprocs=1
startsecs=3
autostart=true
autorestart=true
priority=1
user=root
redirect_stderr=true
stdout_logfile=%(ENV_SYNCTHING_PATH)s/syncthing.log
stdout_logfile_maxbytes=10MB

[program:stdiscosrv]
directory=%(ENV_STDISCOSRV_PATH)s
command=/bin/entrypoint.sh /bin/stdiscosrv
process_name=%(program_name)s
numprocs=1
startsecs=3
autostart=true
autorestart=true
priority=1
user=root
redirect_stderr=true
stdout_logfile=%(ENV_STDISCOSRV_PATH)s/stdiscosrv.log
stdout_logfile_maxbytes=10MB

[program:strelaysrv]
directory=%(ENV_STRELAYSRV_PATH)s
command=/bin/entrypoint.sh /bin/strelaysrv -pools=
process_name=%(program_name)s
numprocs=1
startsecs=3
autostart=true
autorestart=true
priority=1
user=root
redirect_stderr=true
stdout_logfile=%(ENV_STRELAYSRV_PATH)s/strelaysrv.log
stdout_logfile_maxbytes=10MB
