# xray/supervisor.d/xray.ini
[program:xray]
command=/usr/bin/xray -config /etc/xray/config.json
process_name=%(program_name)s
numprocs=1
startsecs=3
autostart=true
autorestart=true
priority=999
user=root
umask=0000
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0