# crazy/supervisor.d/crazy.ini
[program:boot]
command=/bin/sh /etc/crazy/boot.sh
priority=1
startsecs=0
startretries=1
autostart=true
autorestart=false
user=root
umask=0000
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:caddy]
command=/usr/bin/caddy run --config /etc/crazy/Caddyfile --adapter caddyfile
process_name=%(program_name)s
numprocs=1
startsecs=3
autostart=false
autorestart=true
user=root
umask=0000
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:v2ray]
command=/bin/sh /etc/crazy/v2ray.sh
process_name=%(program_name)s
numprocs=1
startsecs=3
autostart=false
autorestart=true
user=root
umask=0000
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:xray]
command=/bin/sh /etc/crazy/xray.sh
process_name=%(program_name)s
numprocs=1
startsecs=3
autostart=false
autorestart=true
user=root
umask=0000
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:ssserver]
command=/bin/sh /etc/crazy/ssserver.sh
process_name=%(program_name)s
numprocs=1
startsecs=3
autostart=false
autorestart=true
user=root
umask=0000
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
