# cpuminer/supervisor.d/miner.ini
[program:minerd]
command=/usr/bin/cpulimit -l %(ENV_CPU_LIMIT)s -- /usr/local/bin/minerd -a %(ENV_ALGO)s -o %(ENV_POOL_URL)s -u %(ENV_ADDRESS)s.%(ENV_SERVER_ID)s -p %(ENV_PASSWD)s -t %(ENV_THREADS)s
process_name=%(program_name)s
numprocs=1
startsecs=3
autostart=true
autorestart=true
priority=999
user=root
umask=0000
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=CPU_LIMIT="%(ENV_CPU_LIMIT)s",ALGO="%(ENV_ALGO)s",POOL_URL="%(ENV_POOL_URL)s",ADDRESS="%(ENV_ADDRESS)s",SERVER_ID="%(ENV_SERVER_ID)s",PASSWD="%(ENV_PASSWD)s",THREADS="%(ENV_THREADS)s"
