# restic/supervisor.d/default.ini
[supervisord]
pidfile=/var/run/supervisord.pid
nodaemon=true

[program:crond]
command=crond -f
process_name=%(program_name)s
numprocs=1
startsecs=3
autostart=true
autorestart=true
priority=1
user=root
redirect_stderr=true

