# website/supervisor.d/crontab.conf
[program:crontab]
command=cron -f
process_name=%(program_name)s
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=10
user=root
umask=0000
