# website/supervisor.d/init.conf
[program:init]
command=/bin/bash -c 'mkdir -p /data/storage; chmod -R 0777 /data/storage;'
priority=1
autorestart=false
startretries=0
startsecs=0
